# inst/systemd/dshpc-worker.service
[Unit]
Description=dsHPC Worker Daemon
After=network.target
Documentation=https://github.com/isglobal-brge/dsFlower-framework

[Service]
Type=simple
User=rock
Group=rock
Environment=DSHPC_HOME=/srv/dshpc
ExecStart=/usr/bin/Rscript /usr/local/lib/R/site-library/dsHPC/worker/main.R ${DSHPC_HOME}
WorkingDirectory=/srv/dshpc
Restart=always
RestartSec=5
StandardOutput=append:/srv/dshpc/worker.log
StandardError=append:/srv/dshpc/worker.log

# Resource limits
LimitNOFILE=4096
MemoryMax=2G

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/srv/dshpc

[Install]
WantedBy=multi-user.target
