# inst/systemd/dsjobs-worker.service
[Unit]
Description=dsJobs Worker Daemon
After=network.target
Documentation=https://github.com/isglobal-brge/dsFlower-framework

[Service]
Type=simple
User=rock
Group=rock
Environment=DSJOBS_HOME=/srv/dsjobs
ExecStart=/usr/bin/Rscript /usr/local/lib/R/site-library/dsJobs/worker/main.R ${DSJOBS_HOME}
WorkingDirectory=/srv/dsjobs
Restart=always
RestartSec=5
StandardOutput=append:/srv/dsjobs/worker.log
StandardError=append:/srv/dsjobs/worker.log

# Resource limits
LimitNOFILE=4096
MemoryMax=2G

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/srv/dsjobs

[Install]
WantedBy=multi-user.target
